# .crontab
@hourly ~/.brew/export > /dev/null 2>&1
